# weather-kiosk.service
[Unit]
Description=Pi Weather Dashboard kiosk (Chromium)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Wayland: do NOT set WAYLAND_DISPLAY/XDG_RUNTIME_DIR here. labwc-autostart must
# run first and call: systemctl --user import-environment WAYLAND_DISPLAY XDG_RUNTIME_DIR
# so this service inherits the real session's display. Otherwise Chromium fails with
# "Failed to connect to Wayland display".
WorkingDirectory=%h/pi-weather-dashboard
ExecStart=%h/pi-weather-dashboard/kiosk.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
